# home-dir/app/system/runit/dotfiles-stow/run
#!/bin/env bash
cd ../../../../../dotfiles-stow/
stow -t $HOME/ .
sleep 10s
