# home-dir/app/system/runit/dotfiles-stow/run
#!/bin/env bash
(cd ../../floppy-void/dotfiles-stow/ && stow -t $HOME/ .)
sleep 3s
