# etc/sv/runsvdir-simon/run
#!/bin/sh

export USER="simon"
export HOME="/home/simon/"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/app/system/runit/"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
